# Lab1/felskit-lab1-archive/overlay/etc/network/if-up.d/02-ntpdate
#!/bin/sh
ntpdate pl.pool.ntp.org
exit 0
